# backend/deploy/raspberry-pi/northax-api.service
[Unit]
Description=NorthAx API (FastAPI/uvicorn)
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target postgresql.service redis-server.service

[Service]
Type=simple
User=@APP_USER@
WorkingDirectory=@BACKEND_DIR@
EnvironmentFile=@BACKEND_DIR@/.env
# IPv4 (0.0.0.0). uvicorn's `--host ::` binds IPv6-ONLY here (not dual-stack),
# so we serve IPv4 and keep mDNS IPv4-only (avahi use-ipv6=no, see README) so
# the iOS app reaches rafaelpereira.local over IPv4 rather than a global AAAA.
ExecStart=@BACKEND_DIR@/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=3
# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
