# backend/deploy/raspberry-pi/northax-api.service
[Unit]
Description=NorthAx API (FastAPI/uvicorn)
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target postgresql.service redis-server.service

[Service]
Type=simple
User=@APP_USER@
WorkingDirectory=@BACKEND_DIR@
EnvironmentFile=@BACKEND_DIR@/.env
# Bind dual-stack (::) so LAN clients reaching the Pi over IPv6 (e.g. the iOS
# app resolving rafaelpereira.local to an AAAA record) connect too; with
# net.ipv6.bindv6only=0 (the Linux default) this one socket also serves IPv4.
ExecStart=@BACKEND_DIR@/.venv/bin/uvicorn app.main:app --host :: --port 8080
Restart=on-failure
RestartSec=3
# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
